# mcp-server-tmp_1rs_h4z/mcp-server.service
[Unit]
Description=Test API MCP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=mcp
Group=mcp
WorkingDirectory=/opt/test-api
Environment=PATH=/opt/test-api/venv/bin
ExecStart=/opt/test-api/venv/bin/python server.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/test-api/data

# Environment variables
Environment=MCP_HOST=0.0.0.0
Environment=MCP_PORT=8080

[Install]
WantedBy=multi-user.target
